# common/development/bash/profile.d/atom
# Configuration file for atom
hash atom 2> /dev/null && export ATOM_HOME="$XDG_CONFIG_HOME/atom/"
