# common/development/bash/profile.d/npm
# Configuration file for nodejs' npm
if hash npm 2>/dev/null
then
	export NPM_CONFIG_USERCONFIG="$XDG_CONFIG_HOME/npm/npmrc"
fi
