# base/bash/profile.d/npm.sh
# Configuration file for nodejs' npm
if hash npm 2>/dev/null
then
	export PATH="${HOME}/.local/lib/nodejs/bin:${PATH}"
	export NPM_CONFIG_USERCONFIG="$XDG_CONFIG_HOME/npm/npmrc"
fi
