# base/bash/profile.d/rust.sh
#!/bin/bash
# Source default install directory for rustup.rs and cargo
[[ -f "$HOME/.cargo/env" ]] && source $HOME/.cargo/env
