# base/bash/profile.d/rust
# Source default install directory for rustup.rs and cargo
if [ -f "$HOME/.cargo/env" ]; then
    source $HOME/.cargo/env
fi
