# common/graphical/bash/profile.d/firefox
# Faster Firefox Javascript performance for Nightly, Aurora, and Beta
# See http://forums.mozillazine.org/viewtopic.php?f=23&t=1948595&p=13825877&hilit=environment#p13825877
export JSGC_DISABLE_POISONING=1

# Use xinput2
# See https://news.ycombinator.com/item?id=18974228
export MOZ_USE_XINPUT2=1

# Use WebRender
export MOZ_WEBRENDER=1

# Native KDE file dialogs for Firefox and Thunderbird
# See https://www.reddit.com/r/linux/comments/a59ip5/firefox_640_released/eblfwcb/
# Causes crashes w/ Thunderbird <69. Can set globally once an updated
# Thunderbird has been released
#export GTK_USE_PORTAL=1
