# common/torrent/systemd/user/transmission-update-blocklists.timer
[Unit]
Description=Update tranmission's IP blocklists weekly

[Timer]
OnCalendar=weekly
Persistent=true
RandomizedDelaySec=10h

[Install]
WantedBy=timers.target
